# raspi/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
priority=998

[program:cron]
command=/sbin/cron
process_name=%(program_name)s
autorestart=false
startsecs=0
startretries=0